# rel/debian/prerm
#!/bin/sh

if which invoke-rc.d >/dev/null 2>&1; then
    invoke-rc.d eadc-hub stop
else
    /etc/init.d/eadc-hub stop
fi
